# server/config/supervisor/mongodb.conf
[program:mongodb]
command=/opt/mongodb/bin/mongod --dbpath /var/local/mongodb --rest
directory=/opt/mongodb
user=root
autorestart=true
